# tyaira.service
[Unit]
Description=Tyaira Driver Testing Platform
After=network.target postgresql.service

[Service]
Type=notify
User=ubuntu
Group=www-data
WorkingDirectory=/home/ubuntu/Tyaira-Provisonal-Project
Environment="PATH=/home/ubuntu/Tyaira-Provisonal-Project/venv/bin"
Environment="FLASK_ENV=production"

# Load environment variables from .env file
EnvironmentFile=/home/ubuntu/Tyaira-Provisonal-Project/.env

# Gunicorn configuration
ExecStart=/home/ubuntu/Tyaira-Provisonal-Project/venv/bin/gunicorn \
    --workers 3 \
    --worker-class sync \
    --bind unix:/home/ubuntu/Tyaira-Provisonal-Project/tyaira.sock \
    --umask 007 \
    --access-logfile /var/log/tyaira/access.log \
    --error-logfile /var/log/tyaira/error.log \
    --log-level info \
    wsgi:app

# Restart policy
Restart=always
RestartSec=10

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
